# docs/crostini-mpris-bridge.service
[Unit]
Description=Crostini MPRIS Bridge (Rust)
After=default.target

[Service]
Type=simple
ExecStart=%h/.local/bin/crostini-mpris-bridge
Restart=on-failure
RestartSec=2
Environment=BIND_ADDR=0.0.0.0:5000
Environment=POLL_INTERVAL_MS=500

[Install]
WantedBy=default.target
